# platform/rpi-3/files/etc/init.d/display-off
#!/bin/sh /etc/rc.common

START=50

start() {
	tvservice -o | logger -p info -t rpi
}

