# platform/xunlong_orangepi-zero/files/etc/uci-defaults/70_orangepi-zero
#
# orangepi-zero - configuration
#

uci -q batch << EOF
set system.status_led=led
set system.status_led.sysfs=orangepi:red:status
set system.status_led.trigger=activity
commit system
set network.globals.packet_steering=1
commit network
set irqbalance.irqbalance.enabled=1
commit irqbalance
EOF

