# platform/friendlyarm_nanopi-neo/files/etc/init.d/home-dir
#!/bin/sh /etc/rc.common

START=40

HDIR=/var/run/admin

start() {
    mkdir -p "$HDIR"
    chown -R admin:admin "$HDIR"
    chmod 750 "$HDIR"
}
