# platform/olimex_a10-olinuxino-lime/files/etc/uci-defaults/99-a10-olinuxino-lime
#
# A10 oLinuXini Lime - configuration
#

uci -q batch << EOF
set system.usr_led=led
set system.usr_led.sysfs=a10-olinuxino-lime:green:usr
set system.usr_led.trigger=activity
commit system
EOF

